# dot_config/systemd/user/flatpak-install.service
[Unit]
Description=Install Flatpak applications from Flatpakfile

[Service]
ExecStart=sh -c "xargs flatpak install < ~/.Flatpakfile"
